# neighborly.ini
[program:neighborly]
process_name=%(program_name)s_%(process_num)02d
command=/root/.nvm/versions/node/v15.2.1/bin/npm run prod
directory=/var/www/neighborly.tools
autostart=true
autorestart=true
user=ec2-user
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/neighborly.tools/neighborly.log
stopwaitsecs=3600